# systemd/netaudio.service
[Unit]
Description=netaudio daemon for Dante network audio device control
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=/usr/bin/netaudio server start
ExecStop=/usr/bin/netaudio server stop
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
